# klishd.service
[Unit]
Description=Klish CLI Daemon for VPP
After=vpp.service network.target
Wants=vpp.service

[Service]
Type=simple
ExecStart=/usr/local/bin/klishd-wrapper
ExecStop=/usr/bin/pkill -x klishd
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
