# klishd.service
[Unit]
Description=Klish CLI Daemon for VPP
After=vpp.service network.target
Requires=vpp.service

[Service]
Type=simple
ExecStartPre=/bin/sleep 3
ExecStart=/usr/local/bin/klishd
ExecStartPost=/bin/bash -c 'sleep 2 && [ -f /etc/vpp/klish-startup.conf ] && vppctl exec /etc/vpp/klish-startup.conf || true'
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
