# klishd.service
[Unit]
Description=Klish CLI Daemon for VPP
After=vpp.service network.target
Requires=vpp.service

[Service]
Type=simple
ExecStartPre=/bin/sleep 2
ExecStartPre=/usr/local/bin/vpp-load-config
ExecStart=/usr/local/bin/klishd
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
